# data/systemd-units/i3-session.target
[Unit]
Description=i3-session
BindsTo=i3-session.service
Wants=graphical-session.target
